# scripts/agent-event-bus.service
[Unit]
Description=Agent Event Bus - Cross-session MCP communication
After=network.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
Environment=PYTHONPATH=__PROJECT_DIR__/src
Environment=AGENT_EVENT_BUS_ICON=__PROJECT_DIR__/assets/icon.png
ExecStart=__VENV_PYTHON__ -m agent_event_bus.server
Restart=always
RestartSec=5
StandardOutput=append:__HOME__/.claude/contrib/agent-event-bus/agent-event-bus.log
StandardError=append:__HOME__/.claude/contrib/agent-event-bus/agent-event-bus.err

[Install]
WantedBy=default.target
