# scripts/com.evansenter.agent-event-bus.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.evansenter.agent-event-bus</string>

    <key>ProgramArguments</key>
    <array>
        <string>__VENV_PYTHON__</string>
        <string>-m</string>
        <string>agent_event_bus.server</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>__PROJECT_DIR__/src</string>
        <key>AGENT_EVENT_BUS_ICON</key>
        <string>__PROJECT_DIR__/assets/icon.png</string>
        <!-- Multi-machine: Set HOST to 0.0.0.0 to accept remote connections -->
        <!-- <key>HOST</key> -->
        <!-- <string>0.0.0.0</string> -->
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!-- Note: launchd truncates on restart; Python logger handles append to .log file -->
    <key>StandardOutPath</key>
    <string>__HOME__/.claude/contrib/agent-event-bus/agent-event-bus.stdout</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.claude/contrib/agent-event-bus/agent-event-bus.err</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
